# supervisord.conf
; root-supervisord для системных сервисов (gpu-broker).
; запуск:  sudo /usr/local/bin/supervisord -c /etc/gpu-broker/supervisord.conf
; статус:  sudo /usr/local/bin/supervisorctl -c /etc/gpu-broker/supervisord.conf status

[unix_http_server]
file=/run/gpu-broker/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/gpu-broker/supervisord.log
pidfile=/run/gpu-broker/supervisord.pid
childlogdir=/var/log/gpu-broker
nodaemon=false
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/gpu-broker/supervisor.sock

[program:gpu-brokerd]
command=/usr/bin/python3 /usr/local/lib/gpu-broker/brokerd.py
autostart=true
autorestart=true
startsecs=3
stopsignal=TERM
stdout_logfile=/var/log/gpu-broker/brokerd.out.log
stderr_logfile=/var/log/gpu-broker/brokerd.err.log
